# list/waydroid.list
/usr/lib/libgbinder.so.1.1.34
/usr/lib/libgbinder.so.1.1
/usr/lib/libgbinder.so.1
/usr/lib/libgbinder.so
/usr/lib/libglibutil.so.1.0.71
/usr/lib/libglibutil.so.1.0
/usr/lib/libglibutil.so.1
/usr/lib/libglibutil.so
/usr/include/gbinder
/usr/include/gutil
/usr/lib/pkgconfig/libgbinder.pc
/usr/lib/pkgconfig/libglibutil.pc
/usr/lib/waydroid
/usr/bin/waydroid
/usr/share/dbus-1/system.d/id.waydro.Container.conf
/usr/share/dbus-1/system.d/com.openfde.Infra.conf
/usr/share/polkit-1/actions/id.waydro.Container.policy
/usr/share/dbus-1/system-services/id.waydro.Container.service
/usr/share/dbus-1/system-services/com.openfde.Infra.service
/usr/lib/systemd/system/waydroid-container.service
/usr/lib/systemd/system/openfde-infra.service
/etc/udev/rules.d/99-hide-android-apex-loops.rules
