# list/waydroid.list
/usr/lib/libgbinder.so.1.1.34
/usr/lib/libgbinder.so.1.1
/usr/lib/libgbinder.so.1
/usr/lib/libglibutil.so.1.0.71
/usr/lib/libglibutil.so.1.0
/usr/lib/libglibutil.so.1
/usr/lib/waydroid
/usr/bin/waydroid
/usr/share/dbus-1/system.d/id.waydro.Container.conf
/usr/share/polkit-1/actions/id.waydro.Container.policy
/usr/share/dbus-1/system-services/id.waydro.Container.service
/usr/lib/systemd/system/waydroid-container.service
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-6.so
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-6.so.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-6.so.0.0.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/Cogl-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/Cogl-6.typelib
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-path-6.so
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-path-6.so.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-path-6.so.0.0.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-pango-6.so
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-pango-6.so.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-cogl-pango-6.so.0.0.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/CoglPango-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/CoglPango-6.typelib
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-clutter-6.so
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-clutter-6.so.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/libmutter-clutter-6.so.0.0.0
/usr/local/lib/aarch64-linux-gnu/mutter-6/Clutter-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/Clutter-6.typelib
/usr/local/lib/aarch64-linux-gnu/mutter-6/Cally-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/Cally-6.typelib
/usr/local/lib/aarch64-linux-gnu/mutter-6/ClutterX11-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/ClutterX11-6.typelib
/usr/local/lib/aarch64-linux-gnu/libmutter-6.so.0.0.0
/usr/local/lib/aarch64-linux-gnu/libmutter-6.so.0
/usr/local/lib/aarch64-linux-gnu/libmutter-6.so
/usr/local/bin/mutter
/usr/local/libexec/mutter-restart-helper
/usr/local/lib/aarch64-linux-gnu/mutter-6/Meta-6.gir
/usr/local/lib/aarch64-linux-gnu/mutter-6/Meta-6.typelib
/usr/local/lib/aarch64-linux-gnu/mutter-6/plugins/libdefault.so
/usr/local/lib/aarch64-linux-gnu/pkgconfig/mutter-cogl-6.pc
/usr/local/lib/aarch64-linux-gnu/pkgconfig/mutter-cogl-path-6.pc
/usr/local/lib/aarch64-linux-gnu/pkgconfig/mutter-cogl-pango-6.pc
/usr/local/lib/aarch64-linux-gnu/pkgconfig/mutter-clutter-6.pc
/usr/local/lib/aarch64-linux-gnu/pkgconfig/mutter-clutter-x11-6.pc
/usr/local/share/gnome-control-center/keybindings/50-mutter-navigation.xml
/usr/local/share/gnome-control-center/keybindings/50-mutter-system.xml
/usr/local/share/gnome-control-center/keybindings/50-mutter-windows.xml
/usr/local/share/gnome-control-center/keybindings/50-mutter-wayland.xml
/usr/local/share/glib-2.0/schemas/org.gnome.mutter.gschema.xml
/usr/local/share/glib-2.0/schemas/org.gnome.mutter.wayland.gschema.xml
/usr/local/share/GConf/gsettings/mutter-schemas.convert
/usr/local/lib/aarch64-linux-gnu/pkgconfig/libmutter-6.pc
/usr/local/lib/python3.8/dist-packages/gbinder.cpython-38-aarch64-linux-gnu.so
